# src/apps/logger/i01logger.socket
[Socket]
ListenStream=/run/i01logger.sock

[Install]
WantedBy=sockets.target

